# systemd/omarchy-i18n-check.timer
[Unit]
Description=Daily check for new untranslated Omarchy strings

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
